# contrib/init/dogecoind.conf
description "Radio Coin Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env RADCCOIND_BIN="/usr/bin/dogecoind"
env RADCCOIND_USER="dogecoin"
env RADCCOIND_GROUP="dogecoin"
env RADCCOIND_PIDDIR="/var/run/dogecoind"
# upstart can't handle variables constructed with other variables
env RADCCOIND_PIDFILE="/var/run/dogecoind/dogecoind.pid"
env RADCCOIND_CONFIGFILE="/etc/dogecoin/dogecoin.conf"
env RADCCOIND_DATADIR="/var/lib/dogecoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # dogecoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$RADCCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run dogecoind."
        echo "The setting must appear in $RADCCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Radiocoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$RADCCOIND_PIDDIR"
    chmod 0755 "$RADCCOIND_PIDDIR"
    chown $RADCCOIND_USER:$RADCCOIND_GROUP "$RADCCOIND_PIDDIR"
    chown $RADCCOIND_USER:$RADCCOIND_GROUP "$RADCCOIND_CONFIGFILE"
    chmod 0660 "$RADCCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$RADCCOIND_PIDFILE" \
    --chuid $RADCCOIND_USER:$RADCCOIND_GROUP \
    --exec "$RADCCOIND_BIN" \
    -- \
    -pid="$RADCCOIND_PIDFILE" \
    -conf="$RADCCOIND_CONFIGFILE" \
    -datadir="$RADCCOIND_DATADIR" \
    -disablewallet \
    -daemon

